# root/etc/s6-overlay/s6-rc.d/svc-semaphore/run
#!/command/with-contenv bash
exec s6-setuidgid semaphore /usr/local/bin/semaphore-wrapper /usr/local/bin/semaphore server --config /etc/semaphore/config.json